# pi/testate-observe.timer
[Unit]
Description=Daily autoportrait at a random hour

[Timer]
OnCalendar=*-*-* 09:00
RandomizedDelaySec=6h
Persistent=false

[Install]
WantedBy=timers.target
