# pi/testate-twilight.timer
[Unit]
Description=Twilight ritual before the Diem epoch turns at midnight UTC

[Timer]
OnCalendar=*-*-* 21:30:00 UTC
Persistent=false

[Install]
WantedBy=timers.target
